# extras/deb/DEBIAN/postinst
#!/bin/bash

set -e

systemctl daemon-reload

systemctl enable hashserver

exit 0
