# extras/deb/DEBIAN/prerm
#!/bin/bash

set -e

systemctl stop hashserver
systemctl disable hashserver

exit 0
