# extras/deb/DEBIAN/prerm
#!/bin/bash

set -e

systemctl disable hashserver

exit 0
